# redis/crontab
# Weekly Redis backup - Daily stats report is at 11am
0 7 * * 1 bensley /opt/dnas/redis/backup_redis.sh > /tmp/backup_redis.log 2>&1

